# supervisord.conf
[supervisord]
nodaemon=true

[program:startscript]
command=/app/start.sh
autorestart=true
startsecs=2
